# Code-and-environment/ESP-Localization/etc/profile.d/esp32_toolchain.sh
# This file was created by ESP-IDF windows_install_prerequisites.sh
# and will be overwritten if that script is run again.
export PATH="$PATH:/opt/xtensa-esp32-elf/bin"
