# Code-and-environment/ESP-Localization/etc/profile.d/export_idf_path.sh
export IDF_PATH="F:/ESP32/msys32/home/pc/esp/esp-idf"