# Code-and-environment/ESP-Localization/etc/profile.d/export_mdf_path.sh
export MDF_PATH="F:/ESP32/msys32/home/pc/esp/esp-mdf"